# scripts/com.pauldesai.alignment.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pauldesai.alignment.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/mirror-admin/repos/truth-first-beacon/scripts/daily_alignment.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mirror-admin/repos/truth-first-beacon</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mirror-admin/.mirrordna/sandbox/ag/alignment_daily.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mirror-admin/.mirrordna/sandbox/ag/alignment_daily_error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
